# sonos-grouper/sonos-grouper.service
[Unit]
Description=Sonos Grouper Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/marienvanoverbeek/git/inkplate/sonos-grouper
EnvironmentFile=/home/marienvanoverbeek/git/inkplate/sonos-grouper/.env
ExecStart=/home/marienvanoverbeek/git/inkplate/sonos-grouper/.venv/bin/python -m sonos_grouper
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
